# cw3/ex2/tvsctld.service
[Unit]
Description=Tvsapp service
Requires=tvsctld.socket
After=syslog.service 

[Service]
Type=simple
ExecStart=/opt/isel/tvs/tvsctld/tvsctld
User=root
Group=tvsgrp
Restart=on-failure

RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=multi-user.target
Also=tvsctld.socket
